# images/microvm/rootfs/etc/systemd/system/thaw-agent.service
[Unit]
Description=Thaw Agent - MicroVM initialization and GitHub runner setup
After=systemd-networkd.service
Wants=systemd-networkd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/thaw-agent --runner-dir=/home/runner --runner-user=runner
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=LOG_LEVEL=info

[Install]
WantedBy=multi-user.target
