# images/microvm/rootfs/etc/systemd/system/thaw-agent.service
[Unit]
Description=Thaw Agent - MicroVM initialization and GitHub runner setup
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/thaw-agent --runner-dir=/home/gleanuser --runner-user=gleanuser
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=LOG_LEVEL=info

[Install]
WantedBy=multi-user.target
